# systemd/cybersentineldlp-dashboard.service
[Unit]
Description=CyberSentinel DLP Dashboard (Next.js)
Documentation=https://github.com/your-org/cybersentineldlp-dlp
After=network.target cybersentineldlp-server.service
Wants=cybersentineldlp-server.service

[Service]
Type=simple
User=dlp
Group=dlp
WorkingDirectory=/opt/cybersentineldlp-dlp/dashboard
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="NODE_ENV=production"
EnvironmentFile=/opt/cybersentineldlp-dlp/.env
ExecStart=/usr/bin/npm run start -- -H 0.0.0.0

# Performance tuning
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/cybersentineldlp-dlp/dashboard

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cybersentineldlp-dashboard

[Install]
WantedBy=multi-user.target
